# templates/amazon/supervisor.init.erb
#! /bin/sh
#
# skeleton  example file to build /etc/init.d/ scripts.
#       This file should be used to construct scripts for /etc/init.d.
#
#       Written by Miquel van Smoorenburg <miquels@cistron.nl>.
#       Modified for Debian
#       by Ian Murdock <imurdock@gnu.ai.mit.edu>.
#               Further changes by Javier Fernandez-Sanguino <jfs@debian.org>
#
# Version:  @(#)skeleton  1.9  26-Feb-2001  miquels@cistron.nl
#
### BEGIN INIT INFO
# Provides:          supervisor
# Required-Start:    $remote_fs $network $named
# Required-Stop:     $remote_fs $network $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop supervisor
# Description:       Start/stop supervisor daemon and its configured
#                    subprocesses.
### END INIT INFO

. /etc/init.d/functions

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON="$(which supervisord)"
NAME=supervisord
DESC=supervisor

test -n "$DAEMON" || exit 0
test -x $DAEMON || exit 0

LOGDIR=/var/log/supervisor
PIDFILE=/var/run/$NAME.pid
DODTIME=5                   # Time to wait for the server to die, in seconds
                            # If this value is set too low you might not
                            # let some servers to die gracefully and
                            # 'restart' will not work

# Include supervisor defaults if available
if [ -f /etc/default/supervisor ] ; then
    . /etc/default/supervisor
fi

start()
{
    echo -n "Starting $DESC: "
    daemon $DAEMON --pidfile $PIDFILE
    RETVAL=$?
    if [ $RETVAL -eq 0 ] ; then
        echo "$NAME."
    else
        echo " ERROR."
    fi
}

stop()
{
    echo -n "Stopping $DESC: "
    echo
    [ -f $PIDFILE ] && killproc $NAME || success "$NAME."
}


case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    echo -n "Restarting $DESC: "
    stop
    start
    ;;
  status)
    status $DAEMON
    ;;
  *)
    N=/etc/init.d/$NAME
    echo "Usage: $N {start|stop|restart|status}" >&2
    exit 1
    ;;
esac

exit 0
